# main/MIStime.service
[Unit]
Description=Time set up workaround service for MISlocker
Requires=network-online.target
After=multi-user.target
DefaultDependencies=true

[Service]
Type=oneshot
ExecStart=/home/MISlocker/system/main/time_updater_proxy_readonly.sh
Restart=no
# must set user to root to execute all functions and peripherals
User=root

[Install]
WantedBy=sysinit.target multi-user.target
